# infra/oracle/runner/systemd/pgbadger-report.timer
[Unit]
Description=Desk Imperial pgBadger daily timer

[Timer]
OnCalendar=*-*-* 04:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
