# infra/oracle/runner/systemd/restore-check.service
[Unit]
Description=Desk Imperial restore drill
After=docker.service wg-quick@wg0.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/desk-imperial
ExecStart=/usr/bin/env bash /opt/desk-imperial/infra/scripts/oracle-runner-bootstrap.sh restore-check
