# infra/oracle/db/systemd/network-guard.service
[Unit]
Description=Desk Imperial Docker DB network guard
After=docker.service wg-quick@wg0.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/desk-imperial
ExecStart=/usr/bin/env bash /opt/desk-imperial/infra/scripts/oracle-db-network-guard.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
